# automation/51-automedia.rules
ACTION=="change", SUBSYSTEM=="block", RUN+="/opt/fartrippr/automation/rippr_wrappr.sh"